# templates/mongod_upstart.conf.erb
# Ubuntu upstart file at /etc/init/mongod_<%= @mongod_instance %>.conf

start on runlevel [2345]
stop on runlevel [06]

limit nofile 64000 64000

script

  HOMEDIR=<%= @conf['homedir'] %>
  DBUSER=<%= @conf['user'] %>
  CONFIG_FILE=<%= @conf['configfile'] %>


  if [ -f $HOMEDIR/mongod.lock ]; then
          rm $HOMEDIR/mongod.lock
          sudo -u $DBUSER /usr/bin/mongod --config $CONFIG_FILE --repair
  fi
  exec start-stop-daemon --start --quiet --pidfile $HOMEDIR/mongo.pid --chuid $DBUSER --exec  /usr/bin/mongod -- --config $CONFIG_FILE

end script